# nvr/deploy/go2rtc.service
[Unit]
Description=go2rtc RTSP-to-WebRTC relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
ExecStart=/opt/nvr/go2rtc -config /opt/nvr/go2rtc.yaml
WorkingDirectory=/opt/nvr
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
